# var/lib/dpkg/info/tinc.list
/.
/etc
/etc/default
/etc/default/tinc
/etc/init.d
/etc/init.d/tinc
/etc/network
/etc/network/if-post-down.d
/etc/network/if-post-down.d/tinc
/etc/network/if-pre-up.d
/etc/network/if-pre-up.d/tinc
/etc/network/if-up.d
/etc/network/if-up.d/tinc
/etc/tinc
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/tinc.service
/lib/systemd/system/tinc@.service
/usr
/usr/sbin
/usr/sbin/tincd
/usr/share
/usr/share/doc
/usr/share/doc/tinc
/usr/share/doc/tinc/AUTHORS
/usr/share/doc/tinc/NEWS.Debian.gz
/usr/share/doc/tinc/NEWS.gz
/usr/share/doc/tinc/README.Debian
/usr/share/doc/tinc/README.gz
/usr/share/doc/tinc/changelog.Debian.gz
/usr/share/doc/tinc/changelog.gz
/usr/share/doc/tinc/copyright
/usr/share/doc/tinc/tinc.html
/usr/share/doc/tinc/tinc.html/About-us.html
/usr/share/doc/tinc/tinc.html/Authentication-protocol.html
/usr/share/doc/tinc/tinc.html/Authors.html
/usr/share/doc/tinc/tinc.html/Building-and-installing-tinc.html
/usr/share/doc/tinc/tinc.html/Concept-Index.html
/usr/share/doc/tinc/tinc.html/Configuration-files.html
/usr/share/doc/tinc/tinc.html/Configuration-introduction.html
/usr/share/doc/tinc/tinc.html/Configuration-of-Darwin-_0028Mac-OS-X_0029-kernels.html
/usr/share/doc/tinc/tinc.html/Configuration-of-FreeBSD-kernels.html
/usr/share/doc/tinc/tinc.html/Configuration-of-Linux-kernels.html
/usr/share/doc/tinc/tinc.html/Configuration-of-NetBSD-kernels.html
/usr/share/doc/tinc/tinc.html/Configuration-of-OpenBSD-kernels.html
/usr/share/doc/tinc/tinc.html/Configuration-of-Solaris-kernels.html
/usr/share/doc/tinc/tinc.html/Configuration-of-Windows.html
/usr/share/doc/tinc/tinc.html/Configuration.html
/usr/share/doc/tinc/tinc.html/Configuring-the-kernel.html
/usr/share/doc/tinc/tinc.html/Contact-information.html
/usr/share/doc/tinc/tinc.html/Cygwin-_0028Windows_0029-build-environment.html
/usr/share/doc/tinc/tinc.html/Darwin-_0028Mac-OS-X_0029-build-environment.html
/usr/share/doc/tinc/tinc.html/Debug-levels.html
/usr/share/doc/tinc/tinc.html/Device-files.html
/usr/share/doc/tinc/tinc.html/Encryption-of-network-packets.html
/usr/share/doc/tinc/tinc.html/Error-messages.html
/usr/share/doc/tinc/tinc.html/Example-configuration.html
/usr/share/doc/tinc/tinc.html/Generating-keypairs.html
/usr/share/doc/tinc/tinc.html/Host-configuration-variables.html
/usr/share/doc/tinc/tinc.html/How-connections-work.html
/usr/share/doc/tinc/tinc.html/How-to-configure.html
/usr/share/doc/tinc/tinc.html/Installation.html
/usr/share/doc/tinc/tinc.html/Interface-configuration.html
/usr/share/doc/tinc/tinc.html/Introduction.html
/usr/share/doc/tinc/tinc.html/Libraries.html
/usr/share/doc/tinc/tinc.html/LibreSSL_002fOpenSSL.html
/usr/share/doc/tinc/tinc.html/Main-configuration-variables.html
/usr/share/doc/tinc/tinc.html/MinGW-_0028Windows_0029-build-environment.html
/usr/share/doc/tinc/tinc.html/Multiple-networks.html
/usr/share/doc/tinc/tinc.html/Network-interfaces.html
/usr/share/doc/tinc/tinc.html/Other-files.html
/usr/share/doc/tinc/tinc.html/Platform-specific-information.html
/usr/share/doc/tinc/tinc.html/Preparations.html
/usr/share/doc/tinc/tinc.html/Routes.html
/usr/share/doc/tinc/tinc.html/Running-tinc.html
/usr/share/doc/tinc/tinc.html/Runtime-options.html
/usr/share/doc/tinc/tinc.html/Scripts.html
/usr/share/doc/tinc/tinc.html/Security-issues.html
/usr/share/doc/tinc/tinc.html/Security.html
/usr/share/doc/tinc/tinc.html/Sending-bug-reports.html
/usr/share/doc/tinc/tinc.html/Signals.html
/usr/share/doc/tinc/tinc.html/Solving-problems.html
/usr/share/doc/tinc/tinc.html/Supported-platforms.html
/usr/share/doc/tinc/tinc.html/System-files.html
/usr/share/doc/tinc/tinc.html/Technical-information.html
/usr/share/doc/tinc/tinc.html/The-UDP-tunnel.html
/usr/share/doc/tinc/tinc.html/The-connection.html
/usr/share/doc/tinc/tinc.html/The-meta_002dconnection.html
/usr/share/doc/tinc/tinc.html/The-meta_002dprotocol.html
/usr/share/doc/tinc/tinc.html/Virtual-Private-Networks.html
/usr/share/doc/tinc/tinc.html/index.html
/usr/share/doc/tinc/tinc.html/lzo.html
/usr/share/doc/tinc/tinc.html/tinc.html
/usr/share/doc/tinc/tinc.html/zlib.html
/usr/share/doc-base
/usr/share/doc-base/tinc
/usr/share/info
/usr/share/info/tinc.info.gz
/usr/share/locale
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/tinc.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/tincd.8.gz
