# var/lib/dpkg/info/tinc.conffiles
/etc/default/tinc
/etc/init.d/tinc
/etc/network/if-post-down.d/tinc
/etc/network/if-pre-up.d/tinc
/etc/network/if-up.d/tinc
