# etc/network/if-up.d/tinc
#!/bin/sh

set -e

if [ "$METHOD" = loopback -o -n "$IF_TINC_NET" ]; then
	exit 0
fi

invoke-rc.d tinc alarm || exit 0
